# systemd/alerts-collector.service
[Unit]
Description=MBTA Alerts Collector

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /root/mbta-scripts/alerts-collector.py
EnvironmentFile=/etc/mbta-collector.env
WorkingDirectory=/root/mbta-scripts
StandardOutput=append:/root/mbta-scripts/alerts-collector.log
StandardError=append:/root/mbta-scripts/alerts-collector.log