# systemd/alerts-collector.timer
[Unit]
Description=Run Alerts Collector every minute at second 05

[Timer]
OnCalendar=*-*-* *:*:05
AccuracySec=1s
Persistent=true
Unit=alerts-collector.service

[Install]
WantedBy=timers.target